# tg-nodeseek-bot.service
[Unit]
Description=Telegram NodeSeek Verification Bot
After=network.target

[Service]
Type=simple
User=root
WorkingDirectory=/opt/tg-nodeseek-bot
Environment=PYTHONUNBUFFERED=1
EnvironmentFile=/opt/tg-nodeseek-bot/.env
ExecStart=/usr/bin/python3 /opt/tg-nodeseek-bot/main.py
Restart=always
RestartSec=10

[Install]
WantedBy=multi-user.target
